# apps/dashboard/deploy/cloud-threat-modeler-dashboard.service
[Unit]
Description=Cloud Threat Modeler Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=fleet
Group=fleet
WorkingDirectory=/home/fleet/cloud-threat-modeler
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONPATH=/home/fleet/cloud-threat-modeler/src
ExecStart=/home/fleet/cloud-threat-modeler/.venv/bin/uvicorn apps.dashboard.main:app --host 127.0.0.1 --port 8001 --workers 2 --proxy-headers
Restart=on-failure
RestartSec=5
TimeoutStopSec=20
NoNewPrivileges=true
PrivateTmp=true
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
